# judge/deploy/judge-api.service
[Unit]
Description=AI Deep Dive Judge API
After=network.target redis.service

[Service]
Type=simple
User=judge
WorkingDirectory=/opt/ai-deep-dive/judge
EnvironmentFile=/etc/judge/judge.env
ExecStart=/opt/ai-deep-dive/judge/.venv/bin/uvicorn judge.api:app --host 127.0.0.1 --port 8000 --workers 1
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
