# judge/deploy/judge-worker-torch.service
[Unit]
Description=AI Deep Dive Judge Worker (torch)
After=network.target redis.service

[Service]
Type=simple
User=judge
WorkingDirectory=/opt/ai-deep-dive/judge
EnvironmentFile=/etc/judge/judge.env
ExecStart=/opt/ai-deep-dive/judge/.venv/bin/python -m judge.worker --stream queue:torch --group workers-torch --consumer torch-1
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
