# judge/deploy/judge-cleanup.timer
[Unit]
Description=Run judge cleanup daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
